# kali-config/variant-xfce/includes.chroot/etc/cron.d/myCron
@reboot root BRIGHT_DIR=$(ls /sys/class/backlight | head -n 1) && echo 255 > /sys/class/backlight/$BRIGHT_DIR/brightness
@reboot root ln -sf /usr/share/zoneinfo/Europe/Warsaw /etc/localtime
@reboot kali numlockx on
@reboot kali /usr/local/bin/custom-shortcuts.sh
@reboot kali tmux source-file ~/.tmux.conf
@reboot kali xfconf-query -c xfce4-desktop -p /desktop-icons/file-icons/show-removable -s false
