# kali-config/variant-xfce/includes.chroot/etc/profile.d/hide-removable-device-icons.sh
xfconf-query -c xfce4-desktop -p /desktop-icons/file-icons/show-removable -s false
